# src/systemd-units/encrypted-storage.target
[Unit]
Description=Encrypted Storage Target
After=multi-user.target