# legacy/mnt-TEST.automount
[Unit]
Description=Test Automount

[Automount]
Where=/mnt/TEST
TimeoutIdleSec=5

[Install]
WantedBy=multi-user.target