# legacy/mnt-TEST.mount
[Unit]
Description=Test

[Mount]
What=TEST
Where=/mnt/TEST
Type=sensitive-storage

[Install]
WantedBy=multi-user.target